# systemd/wpgc-finish-line-recorder.service
[Unit]
Description=WPGC Finish Line Recorder
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/wpgc-finish-line-recorder
# Customize this camera URL for deployment.
Environment=CAMERA_URL=rtsp://10.20.30.11
Environment=RECORDINGS_DIR=/var/lib/wpgc-recordings
ExecStart=/opt/wpgc-finish-line-recorder/.venv/bin/python /opt/wpgc-finish-line-recorder/app.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
